# src/packaging/bin/init-script/hivemq.service
[Unit]
Description=HiveMQ Community Edition
After=network.target

[Service]
Group=hivemq
User=hivemq
LimitNOFILE=infinity
StandardOutput=null
StandardError=journal
WorkingDirectory=/opt/hivemq
ExecStart=/opt/hivemq/bin/run.sh
KillMode=process
Restart=always
KillSignal=15
SuccessExitStatus=143
TimeoutStopSec=3600

[Install]
WantedBy=multi-user.target
